# systemd/daybreak-dark.service
[Unit]
Description=Daybreak Dark Theme Switcher (Night Mode)

[Service]
Type=oneshot
ExecStart=/usr/bin/env daybreak dark
IOSchedulingClass=idle
